# supervisor/conf.d/prosody.conf
[program:prosody]
command = /usr/bin/prosody
user = prosody
autostart = true