# prosody/cron.hourly/update-roster
#!/bin/bash

python3 /etc/prosody/ldap-roster.py

/usr/bin/supervisorctl restart prosody
